# infra/systemd/xllmapi-platform.service
[Unit]
Description=xllmapi Platform API
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=xllmapi
WorkingDirectory=/opt/xllmapi
Environment=XLLMAPI_ENV=production
Environment=XLLMAPI_DB_DRIVER=postgres
Environment=XLLMAPI_CORS_ORIGINS=https://app.example.com
Environment=XLLMAPI_SESSION_COOKIE_NAME=xllmapi_session
ExecStart=/usr/bin/node apps/platform-api/dist/main.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
